# pkg/debian/DEBIAN/postinst
#!/bin/sh

set -e

echo "Enable and start systemd service..."
deb-systemd-helper enable ads-catcher.service
deb-systemd-invoke start ads-catcher.service

